# ssh_ssl_proxy/scripts/DEBIAN/postrm
#!/bin/sh
 
# Source debconf library.
. /usr/share/debconf/confmodule

/usr/sbin/update-rc.d -f ssh_ssl_proxy remove